# cntools_py3/cndeploy/cndeploy/templates/supervisor.conf
[program:{{PROGRAM}}]
command={{COMMAND}}
driver
autorestart=true
redirect_stderr=true
startretries=100
user=ubuntu
